# out/deploy/reverse-obfuscated-ssh-out.container
[Unit]
Description=Reverse Obfuscated SSH Container - OUT
After=network.target

[Container]
Image=quay.io/luisarizmendi/con-out:latest
ContainerName=reverse-obfuscated-ssh-out
PublishPort=2223:2223
PublishPort=8089:8189
PublishPort=8443:8443
Environment=OUT_PASS=R3dh4t1!supersecure
Environment=OUT_SSH_OBF_KEY=123qwersuper
[Service]
Restart=always

[Install]
WantedBy=default.target